# export/templates/systemd/default/control.target.tmpl
[Unit]
Wants={{ range $value := .processes }}{{ $value }} {{ end }}

[Install]
WantedBy=multi-user.target
