# export/templates/upstart/default/control.conf.tmpl
description {{ .app }}

start on (started networking)
stop on runlevel [06]
