# export/templates/upstart/default/process-type.conf.tmpl
description {{ .app }}-{{ .process_type }}

start on starting {{ .app }}
stop on (runlevel [06] or stopping {{ .app }})
